# var/lib/pacman/local/wicd-1.7.3-1/files
%FILES%
etc/
etc/dbus-1/
etc/dbus-1/system.d/
etc/dbus-1/system.d/wicd.conf
etc/logrotate.d/
etc/logrotate.d/wicd.logrotate
etc/wicd/
etc/wicd/dhclient.conf.template.default
etc/wicd/encryption/
etc/wicd/encryption/templates/
etc/wicd/encryption/templates/active
etc/wicd/encryption/templates/active_wired
etc/wicd/encryption/templates/eap
etc/wicd/encryption/templates/eap-tls
etc/wicd/encryption/templates/leap
etc/wicd/encryption/templates/peap
etc/wicd/encryption/templates/peap-tkip
etc/wicd/encryption/templates/psu
etc/wicd/encryption/templates/ttls
etc/wicd/encryption/templates/wep-hex
etc/wicd/encryption/templates/wep-passphrase
etc/wicd/encryption/templates/wep-shared
etc/wicd/encryption/templates/wired_8021x
etc/wicd/encryption/templates/wpa
etc/wicd/encryption/templates/wpa-peap
etc/wicd/encryption/templates/wpa-psk
etc/wicd/encryption/templates/wpa-psk-hex
etc/wicd/encryption/templates/wpa2-leap
etc/wicd/encryption/templates/wpa2-peap
etc/wicd/scripts/
etc/wicd/scripts/.empty_on_purpose
etc/wicd/scripts/postconnect/
etc/wicd/scripts/postconnect/.empty_on_purpose
etc/wicd/scripts/postdisconnect/
etc/wicd/scripts/postdisconnect/.empty_on_purpose
etc/wicd/scripts/preconnect/
etc/wicd/scripts/preconnect/.empty_on_purpose
etc/wicd/scripts/predisconnect/
etc/wicd/scripts/predisconnect/.empty_on_purpose
usr/
usr/bin/
usr/bin/wicd
usr/bin/wicd-cli
usr/bin/wicd-curses
usr/lib/
usr/lib/pm-utils/
usr/lib/pm-utils/sleep.d/
usr/lib/pm-utils/sleep.d/55wicd
usr/lib/python2.7/
usr/lib/python2.7/site-packages/
usr/lib/python2.7/site-packages/wicd-1.7.3-py2.7.egg-info
usr/lib/python2.7/site-packages/wicd/
usr/lib/python2.7/site-packages/wicd/__init__.py
usr/lib/python2.7/site-packages/wicd/__init__.pyc
usr/lib/python2.7/site-packages/wicd/__init__.pyo
usr/lib/python2.7/site-packages/wicd/backend.py
usr/lib/python2.7/site-packages/wicd/backend.pyc
usr/lib/python2.7/site-packages/wicd/backend.pyo
usr/lib/python2.7/site-packages/wicd/configmanager.py
usr/lib/python2.7/site-packages/wicd/configmanager.pyc
usr/lib/python2.7/site-packages/wicd/configmanager.pyo
usr/lib/python2.7/site-packages/wicd/dbusmanager.py
usr/lib/python2.7/site-packages/wicd/dbusmanager.pyc
usr/lib/python2.7/site-packages/wicd/dbusmanager.pyo
usr/lib/python2.7/site-packages/wicd/logfile.py
usr/lib/python2.7/site-packages/wicd/logfile.pyc
usr/lib/python2.7/site-packages/wicd/logfile.pyo
usr/lib/python2.7/site-packages/wicd/misc.py
usr/lib/python2.7/site-packages/wicd/misc.pyc
usr/lib/python2.7/site-packages/wicd/misc.pyo
usr/lib/python2.7/site-packages/wicd/networking.py
usr/lib/python2.7/site-packages/wicd/networking.pyc
usr/lib/python2.7/site-packages/wicd/networking.pyo
usr/lib/python2.7/site-packages/wicd/translations.py
usr/lib/python2.7/site-packages/wicd/translations.pyc
usr/lib/python2.7/site-packages/wicd/translations.pyo
usr/lib/python2.7/site-packages/wicd/wnettools.py
usr/lib/python2.7/site-packages/wicd/wnettools.pyc
usr/lib/python2.7/site-packages/wicd/wnettools.pyo
usr/lib/python2.7/site-packages/wicd/wpath.py
usr/lib/python2.7/site-packages/wicd/wpath.pyc
usr/lib/python2.7/site-packages/wicd/wpath.pyo
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/wicd.service
usr/lib/wicd/
usr/lib/wicd/__init__.py
usr/lib/wicd/backend.py
usr/lib/wicd/configmanager.py
usr/lib/wicd/dbusmanager.py
usr/lib/wicd/logfile.py
usr/lib/wicd/misc.py
usr/lib/wicd/networking.py
usr/lib/wicd/translations.py
usr/lib/wicd/wnettools.py
usr/lib/wicd/wpath.py
usr/share/
usr/share/dbus-1/
usr/share/dbus-1/system-services/
usr/share/dbus-1/system-services/org.wicd.daemon.service
usr/share/doc/
usr/share/doc/wicd/
usr/share/doc/wicd/AUTHORS
usr/share/doc/wicd/CHANGES
usr/share/doc/wicd/INSTALL
usr/share/doc/wicd/LICENSE
usr/share/doc/wicd/README
usr/share/doc/wicd/README.cli
usr/share/doc/wicd/README.curses
usr/share/locale/
usr/share/locale/ast/
usr/share/locale/ast/LC_MESSAGES/
usr/share/locale/ast/LC_MESSAGES/wicd.mo
usr/share/locale/ca/
usr/share/locale/ca/LC_MESSAGES/
usr/share/locale/ca/LC_MESSAGES/wicd.mo
usr/share/locale/cs/
usr/share/locale/cs/LC_MESSAGES/
usr/share/locale/cs/LC_MESSAGES/wicd.mo
usr/share/locale/da/
usr/share/locale/da/LC_MESSAGES/
usr/share/locale/da/LC_MESSAGES/wicd.mo
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/wicd.mo
usr/share/locale/el/
usr/share/locale/el/LC_MESSAGES/
usr/share/locale/el/LC_MESSAGES/wicd.mo
usr/share/locale/en/
usr/share/locale/en/LC_MESSAGES/
usr/share/locale/en/LC_MESSAGES/wicd.mo
usr/share/locale/en_GB/
usr/share/locale/en_GB/LC_MESSAGES/
usr/share/locale/en_GB/LC_MESSAGES/wicd.mo
usr/share/locale/es/
usr/share/locale/es/LC_MESSAGES/
usr/share/locale/es/LC_MESSAGES/wicd.mo
usr/share/locale/et/
usr/share/locale/et/LC_MESSAGES/
usr/share/locale/et/LC_MESSAGES/wicd.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/wicd.mo
usr/share/locale/fr_CA/
usr/share/locale/fr_CA/LC_MESSAGES/
usr/share/locale/fr_CA/LC_MESSAGES/wicd.mo
usr/share/locale/gl/
usr/share/locale/gl/LC_MESSAGES/
usr/share/locale/gl/LC_MESSAGES/wicd.mo
usr/share/locale/he/
usr/share/locale/he/LC_MESSAGES/
usr/share/locale/he/LC_MESSAGES/wicd.mo
usr/share/locale/hr/
usr/share/locale/hr/LC_MESSAGES/
usr/share/locale/hr/LC_MESSAGES/wicd.mo
usr/share/locale/hu/
usr/share/locale/hu/LC_MESSAGES/
usr/share/locale/hu/LC_MESSAGES/wicd.mo
usr/share/locale/it/
usr/share/locale/it/LC_MESSAGES/
usr/share/locale/it/LC_MESSAGES/wicd.mo
usr/share/locale/ja/
usr/share/locale/ja/LC_MESSAGES/
usr/share/locale/ja/LC_MESSAGES/wicd.mo
usr/share/locale/nl/
usr/share/locale/nl/LC_MESSAGES/
usr/share/locale/nl/LC_MESSAGES/wicd.mo
usr/share/locale/pl/
usr/share/locale/pl/LC_MESSAGES/
usr/share/locale/pl/LC_MESSAGES/wicd.mo
usr/share/locale/pt/
usr/share/locale/pt/LC_MESSAGES/
usr/share/locale/pt/LC_MESSAGES/wicd.mo
usr/share/locale/pt_BR/
usr/share/locale/pt_BR/LC_MESSAGES/
usr/share/locale/pt_BR/LC_MESSAGES/wicd.mo
usr/share/locale/ru/
usr/share/locale/ru/LC_MESSAGES/
usr/share/locale/ru/LC_MESSAGES/wicd.mo
usr/share/locale/si/
usr/share/locale/si/LC_MESSAGES/
usr/share/locale/si/LC_MESSAGES/wicd.mo
usr/share/locale/sk/
usr/share/locale/sk/LC_MESSAGES/
usr/share/locale/sk/LC_MESSAGES/wicd.mo
usr/share/locale/sv/
usr/share/locale/sv/LC_MESSAGES/
usr/share/locale/sv/LC_MESSAGES/wicd.mo
usr/share/locale/tr/
usr/share/locale/tr/LC_MESSAGES/
usr/share/locale/tr/LC_MESSAGES/wicd.mo
usr/share/locale/uk/
usr/share/locale/uk/LC_MESSAGES/
usr/share/locale/uk/LC_MESSAGES/wicd.mo
usr/share/locale/zh_CN/
usr/share/locale/zh_CN/LC_MESSAGES/
usr/share/locale/zh_CN/LC_MESSAGES/wicd.mo
usr/share/locale/zh_TW/
usr/share/locale/zh_TW/LC_MESSAGES/
usr/share/locale/zh_TW/LC_MESSAGES/wicd.mo
usr/share/man/
usr/share/man/man1/
usr/share/man/man1/wicd-client.1.gz
usr/share/man/man5/
usr/share/man/man5/wicd-manager-settings.conf.5.gz
usr/share/man/man5/wicd-wired-settings.conf.5.gz
usr/share/man/man5/wicd-wireless-settings.conf.5.gz
usr/share/man/man8/
usr/share/man/man8/wicd-cli.8.gz
usr/share/man/man8/wicd-curses.8.gz
usr/share/man/man8/wicd.8.gz
usr/share/man/nl/
usr/share/man/nl/man1/
usr/share/man/nl/man1/wicd-client.1.gz
usr/share/man/nl/man5/
usr/share/man/nl/man5/wicd-manager-settings.conf.5.gz
usr/share/man/nl/man5/wicd-wired-settings.conf.5.gz
usr/share/man/nl/man5/wicd-wireless-settings.conf.5.gz
usr/share/man/nl/man8/
usr/share/man/nl/man8/wicd-curses.8.gz
usr/share/man/nl/man8/wicd.8.gz
usr/share/wicd/
usr/share/wicd/backends/
usr/share/wicd/backends/be-external.py
usr/share/wicd/backends/be-ioctl.py
usr/share/wicd/cli/
usr/share/wicd/cli/wicd-cli.py
usr/share/wicd/curses/
usr/share/wicd/curses/configscript_curses.py
usr/share/wicd/curses/curses_misc.py
usr/share/wicd/curses/netentry_curses.py
usr/share/wicd/curses/prefs_curses.py
usr/share/wicd/curses/wicd-curses.py
usr/share/wicd/daemon/
usr/share/wicd/daemon/autoconnect.py
usr/share/wicd/daemon/monitor.py
usr/share/wicd/daemon/suspend.py
usr/share/wicd/daemon/wicd-daemon.py
usr/share/wicd/scripts/
usr/share/wicd/scripts/50-wicd-suspend.sh
usr/share/wicd/scripts/80-wicd-connect.sh
var/
var/lib/
var/lib/wicd/
var/lib/wicd/WHEREAREMYFILES
var/lib/wicd/configurations/
var/lib/wicd/configurations/.empty_on_purpose
var/log/
var/log/wicd/
var/log/wicd/.empty_on_purpose

%BACKUP%
etc/wicd/encryption/templates/active	6d1626431d222eb0012736c42562d22c

